# week.cron
#!/bin/bash

cd /var/www/hearts

source venv/bin/activate

python stats.py
